# tornado/tornado_ubuntu_service.conf
start on runlevel [2345]
stop on runlevel [!2345]
 
respawn
 
 
exec python2 /var/www/app/server.py
